# build/darp-118f291-dirty/darp.service
[Unit]
Description=DARP Cloudflare WARP Client
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/darp connect
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
